# etc/systemd/system/music-sync-worker.timer
[Unit]
Description=Run music-sync worker periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target